# src/etc/init.d/pljobmgr-redhat.in
#!/bin/sh -f
#
# chkconfig: 2345 99 05
# description: Start/Stop the Pipeline (v@VERSION@) Job Server daemon
# processname: pljobmgr
# pidfile:

PLJOBMGR_BIN=@ROOT_INST_PARENT_DIR@/pipeline-latest/Unix-@NATIVE_ARCHTYPE@-Opt/sbin/pljobmgr
test -f $PLJOBMGR_BIN || { echo "$PLJOBMGR_BIN not installed"; 
	if [ "$1" = "stop" ]; then exit 0;
	else exit 5; fi; }

JPS_BIN=@UNIX_JAVA_HOME@/../bin/jps
test -f $JPS_BIN || { echo "$JPS_BIN not installed"; 
	if [ "$1" = "stop" ]; then exit 0;
	else exit 5; fi; }

. /etc/rc.d/init.d/functions

function running {
  [ "$BOOTUP" = "color" ] && $MOVE_TO_COL
  echo -n "["
  [ "$BOOTUP" = "color" ] && $SETCOLOR_SUCCESS
  echo -n $"RUNNING"
  [ "$BOOTUP" = "color" ] && $SETCOLOR_NORMAL
  echo -n "]"
}

function not_running {
  [ "$BOOTUP" = "color" ] && $MOVE_TO_COL
  echo -n "["
  [ "$BOOTUP" = "color" ] && $SETCOLOR_FAILURE
  echo -n $" DEAD "
  [ "$BOOTUP" = "color" ] && $SETCOLOR_NORMAL
  echo -n "]"
}

function checkJobMgr {
    pid=`$JPS_BIN -m | grep JobMgrApp | awk '{print $1}'`
    if [ "$pid"x != x ]
    then
	return 0
    else
	return 1
    fi
}

function startJobMgr {
    if checkJobMgr
    then
	return 0
    else
	nice -20 runuser -s /bin/bash -g @PIPELINE_GROUP@ - @PIPELINE_USER@ -c "$PLJOBMGR_BIN --standard-log-file --fail-fast" >/dev/null 2>&1 &
	sleep 3
    fi
    checkJobMgr
}

function stopJobMgr {
  pid=`$JPS_BIN -m | grep JobMgrApp | awk '{print $1}'`
  if [ "$pid"x != x ]
  then
      kill -9 $pid
      echo_success
  else
      not_running
      return 1
  fi
}


case "$1" in
  start)
    echo -n "Starting Pipeline Job Manager: " 
    if startJobMgr
    then 
	running
    else
	not_running
    fi    
    ;;

  stop)
    echo -n "Stopping Pipeline Job Manager: " 
    stopJobMgr
    ;;

  status)
    echo -n "Checking for Job Manager Mananager: "
    if checkJobMgr
    then 
	running
    else
	not_running
    fi
    ;;

  *)
    echo "usage: $0 {start|stop|staus}"
    ;;
esac
